# rpms/EL9/enduro-sfa/files/enduro-am-worker.service
[Unit]
Description=Enduro

[Service]
Type=simple
ExecStart=/usr/bin/enduro-am-worker -config /etc/enduro/enduro.toml

[Install]
WantedBy=multi-user.target
